# config/supervisord.ini
[program:jupyter]
command=/opt/anaconda/bin/jupyter-notebook --ip 0.0.0.0 --no-browser --notebook-dir=/home/jupyter/notebooks/
directory=/home/jupyter/runtime/
autostart=true
autorestart=true
user=jupyter
stdout_logfile=/home/jupyter/runtime/jupyter.log
stderr_logfile=/home/jupyter/runtime/jupyter-errors.log

