# config/supervisord.ini
[program:jupyterhub]
command=/usr/bin/jupyterhub --ip 0.0.0.0  --port 8888 -f /etc/jupyterhub_config.py
directory=/var/run/
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/jupyterhub.log
stderr_logfile=/var/log/jupyterhub-errors.log

[program:scrapyd]
command=/usr/bin/scrapyd
directory=/var/run/
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/scrapyd.log
stderr_logfile=/var/log/scrapyd-errors.log

