# config/supervisord.ini
[program:jupyter]
command=/opt/anaconda/bin/jupyter-notebook --ip 0.0.0.0 --no-browser --notebook-dir=/home/runtime/notebooks/
directory=/home/runtime/run/
autostart=true
autorestart=true
user=runtime
stdout_logfile=/home/runtime/run/jupyter.log
stderr_logfile=/home/runtime/run/jupyter-errors.log

